# minidlna-db-autorefresh.timer
[Unit]
Description=Run minidlna-db-autorefresh every day

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
